# software/mother.service
[Unit]
Description=NOSTROMO Terminal
After=multi-user.target

[Service]
Type=simple
User=pi
Environment=ANTHROPIC_API_KEY=your-key-here
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=SDL_AUDIODRIVER=alsa
Environment=AUDIODEV=plughw:2
WorkingDirectory=/home/pi/mother
ExecStart=/usr/bin/python3 /home/pi/mother/main.py --fullscreen
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
